# rootfs/etc/systemd/system/nvleaf-bluetooth.service
[Unit]
Description=NvLeaf Bluetooth Setup Service
Before=bluetooth.service
RequiredBy=bluetooth.service

[Service]
ExecStart=/opt/nvleaf/bluetooth.sh
Type=simple

[Install]
WantedBy=multi-user.target
